# systemd/devicebox.service
[Unit]
Description=DeviceBox Next.js Application
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/devicebox
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

